# etc/rc.d/rc.local

#!/bin/bash
#
# Inicializar a conexão PPPoE na inicialização
/usr/sbin/pppd call pppoe &

